# crates/systemd-parser/tests/success_units/lvm2-activation.service
# Automatically generated by lvm2-activation-generator.
#
# This unit is responsible for direct activation of LVM2 logical volumes
# if lvmetad daemon is not used (global/use_lvmetad=0 lvm.conf setting),
# hence volume autoactivation is not applicable.
# Direct LVM2 activation requires udev to be settled!

[Unit]
Description=Activation of LVM2 logical volumes
Documentation=man:lvm2-activation-generator(8)
SourcePath=/etc/lvm/lvm.conf
DefaultDependencies=no
After= lvm2-activation-early.service cryptsetup.target
Before=local-fs.target shutdown.target
Wants=systemd-udev-settle.service

[Service]
ExecStart=/sbin/lvm vgchange -aay --sysinit --ignoreskippedcluster
Type=oneshot
