# crates/systemd-parser/tests/success_units/ntpd.service
[Unit]
Description=Network Time Service
After=ntpdate.service sntp.service
Conflicts=systemd-timesyncd.service

[Service]
ExecStart=/usr/sbin/ntpd -g -n -f /var/lib/ntp/ntp.drift -u ntp:ntp
PrivateTmp=true
Restart=always

[Install]
WantedBy=multi-user.target
