# crates/systemd-parser/tests/success_units/uuidd.socket
[Unit]
Description=UUID daemon activation socket

[Socket]
ListenStream=/var/lib/uuidd/request

[Install]
WantedBy=sockets.target
